# run/discord-rss-webhook/discord-rss-webhook.service
[Unit]
Description=Discord RSS Webhook
After=network.target

[Service]
Type=oneshot
User=exedev
WorkingDirectory=/home/exedev
EnvironmentFile=/home/exedev/discord-rss-webhook.env
ExecStart=/home/exedev/bin/discord-rss-webhook
